# supervisord.conf
[supervisord]
nodaemon=true


[program:nginx]
command=/usr/servers/nginx/sbin/nginx -g "daemon off;"